# roles/supervisor_at_boot/files/supervisor.conf
description "supervisord"

start on (local-filesystems and net-device-up)
stop on [!12345]

limit nofile 20000

respawn

chdir /tmp

exec /usr/local/bin/supervisord -c /etc/supervisord.conf
